# .local/sv/mpd/run
#!/bin/sh

exec >/dev/null
exec mpd --no-daemon
