# .local/sv/ssh-agent/finish
#!/bin/sh

rm -f $HOME/.ssh/agent
exit 0
