# .local/sv/ssh-agent/run
#!/bin/sh

exec >/dev/null
exec ssh-agent -Da $HOME/.ssh/agent
